# s6-overlay/s6-rc.d/06-librespot/finish
#!/command/with-contenv sh
echo "librespot-api exited with code $?"
exit 0